# chrony-nts/rootfs/etc/s6-overlay/s6-rc.d/init-chrony/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ============================================================================== 
# Home Assistant Add-on: Chrony NTP Server (with NTS)
# Configures chrony before starting the daemon
# ==============================================================================
set -Eeuo pipefail

readonly CONF="/etc/chrony/chrony.conf"
readonly DATA_DIR="/data/chrony"

bashio::log.info "Configuring chrony..."

# Create data directory
mkdir -p "$DATA_DIR"
chmod 0700 "$DATA_DIR" || true

# Get configuration values
MODE="$(bashio::config 'mode')"
POOL="$(bashio::config 'ntp_pool' || true)"
ALLOW_RAW="$(bashio::config 'allow' | tr ',' ' ')"
USE_NTS="$(bashio::config 'use_nts')"
LOCAL_FALLBACK="$(bashio::config 'local_fallback')"
EXTRA_CONFIG="$(bashio::config 'extra_config' || true)"

# Read server list (array)
mapfile -t SERVERS < <(bashio::config 'ntp_server[]' 2>/dev/null || true)

# Build configuration file
{
  echo "# Generated by Home Assistant Chrony NTS add-on"
  echo "driftfile ${DATA_DIR}/chrony.drift"
  echo "makestep 1.0 3"
  echo "rtcsync"

  # ACLs - process the allow list
  for net in ${ALLOW_RAW}; do
    if [ -n "$net" ]; then
      echo "allow ${net}"
    fi
  done

  # Configure time sources based on mode
  if [ "${MODE}" = "pool" ]; then
    if [ -n "${POOL:-}" ]; then
      if bashio::var.true "${USE_NTS}"; then
        echo "pool ${POOL} iburst nts"
      else
        echo "pool ${POOL} iburst"
      fi
    fi
  else
    # Server mode - use individual servers
    for host in "${SERVERS[@]}"; do
      if [ -n "${host:-}" ]; then
        if bashio::var.true "${USE_NTS}"; then
          echo "server ${host} iburst nts"
        else
          echo "server ${host} iburst"
        fi
      fi
    done
  fi

  # Persist NTS cookies if enabled
  if bashio::var.true "${USE_NTS}"; then
    echo "ntsdumpdir ${DATA_DIR}"
  fi

  # Optional local fallback (off by default)
  if bashio::var.true "${LOCAL_FALLBACK}"; then
    echo "local stratum 10"
  fi

  # Advanced user-provided lines (appended verbatim)
  if [ -n "${EXTRA_CONFIG:-}" ]; then
    echo
    echo "# --- extra_config (verbatim) ---"
    echo "${EXTRA_CONFIG}"
  fi
} > "${CONF}"

bashio::log.info "chrony.conf generated successfully at ${CONF}"
